# supervisor.ini
[program:upload]
command=/usr/bin/python3 server.py -b 127.0.0.1:8001 -w 2
numprocs=1
directory=/home/webmaster/app
autostart=true
autorestart=true
startretries=3
user=webmaster

[program:nginx]
command=/usr/sbin/nginx
numprocs=1
autostart=true
autorestart=true
startretries=3
user=webmaster

[group:cdn]
programs=nginx,upload
